# deploy/standalone/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── PostgreSQL ──────────────────────────────────────────────────
[program:postgres]
command=gosu postgres postgres -D /var/lib/postgresql/data -c listen_addresses=127.0.0.1
autostart=true
autorestart=true
priority=1
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Next.js App ─────────────────────────────────────────────────
[program:app]
; DSN is injected by entrypoint.sh via export, but we also set it here
; as a fallback. supervisord inherits the environment from the parent shell.
command=node /app/server.js
autostart=true
autorestart=true
priority=10
; Wait for postgres to be ready before starting the app
startdelay=8
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    NODE_ENV="production",
    PORT="3000",
    ENABLE_RATE_LIMIT="false",
    AUTO_MIGRATE="%(ENV_AUTO_MIGRATE)s",
    ADMIN_TOKEN="%(ENV_ADMIN_TOKEN)s"
